# etc/cron.d/glpi
*/2 * * * * www-data /opt/glpi/glpi-cron.sh >> /dev/null



